# files/init.conf
description "json_template_api"
version "1.0"
author "Rob Wilson"

respawn
start on runlevel [2345]

exec thin start -C /opt/json_template_api/thin.yml start
